# deploy/celery-worker.service
[Unit]
Description=TrackAudit — Celery Worker
After=network.target redis.service
Requires=redis.service

[Service]
Type=forking
User=trackaudit
Group=trackaudit
WorkingDirectory=/trackaudit/backend
EnvironmentFile=/trackaudit/backend/.env
ExecStart=/trackaudit/backend/venv/bin/celery \
    -A trackaudit worker \
    --loglevel=info \
    --concurrency=2 \
    --pidfile=/var/run/celery/worker.pid \
    --logfile=/var/log/celery/worker.log \
    --detach
ExecStop=/trackaudit/backend/venv/bin/celery \
    -A trackaudit control shutdown
PIDFile=/var/run/celery/worker.pid
RuntimeDirectory=celery
RuntimeDirectoryMode=0755
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
